# ops/com.jyufu.rss.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.jyufu.rss.update</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>__REPO_ROOT__/scripts/run-publish-feeds.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/logs/launchd.err.log</string>
  </dict>
</plist>
